# examples/local-runner/filebrowser/monitord
#!/sbin/openrc-run

name="Monitor data dir and change Permissions"
description="Monitor data and change permissions for newly created files"
command="/usr/check_file.sh"
command_background=true
pidfile="/run/monitord.pid"
